# docker/cron/etc/periodic/15min/run-script
#!/bin/sh

echo "Hello cron! >> /var/log/cron.log 2>&1
